# pkg/deb/control/prerm.sh
#!/bin/sh -e
if [ -x "/etc/init.d/disthene" ]; then
	invoke-rc.d disthene stop || exit $?
fi
